# Azure API Example/azureAPI.service
[Unit]
Description=AzureAPI

[Service]
ExecStart=/usr/local/bin/Azure_API_Service
Restart=always
User=root

[Install]
WantedBy=multi-user.target